# roles/metabase/templates/metabase.ini
[program:metabase]
environment=
    PATH="{{JAVA_HOME}}/bin:%(ENV_PATH)s",
    MB_PASSWORD_COMPLEXITY=normal,
    MB_DB_TYPE=h2,
    MB_DB_FILE={{METABASE_DATA_DIR}}/h2.db,
    MB_JETTY_HOST=0.0.0.0,
    MB_JETTY_PORT={{METABASE_PORT}}
command={{JAVA_HOME}}/bin/java -Duser.timezone=Asia/Shanghai -jar metabase-{{METABASE_VERSION}}.jar
stdout_logfile={{METABASE_DATA_DIR}}/logs/metabase.log
directory={{METABASE_HOME}}
redirect_stderr=true