# roles/hive/templates/hive-metastore.ini
[program:hive-metastore]
environment=PATH="{{JAVA_HOME}}/bin:%(ENV_PATH)s",HADOOP_PREFIX="{{HADOOP_HOME}}"
command={{HIVE_HOME}}/bin/hive --service metastore
stdout_logfile={{HIVE_LOG_DIR}}/metastore.log
redirect_stderr=true