# roles/hadoop/templates/services/namenode.ini.j2
[program:namenode]
command={{HADOOP_HOME}}/bin/hdfs namenode
stdout_logfile={{HADOOP_LOGS}}/namenode-service.log
stderr_logfile={{HADOOP_LOGS}}/namenode-error.log