# roles/presto/templates/presto.ini
[program:presto]
environment=PATH="{{JAVA_HOME}}/bin:%(ENV_PATH)s"
command={{PRESTO_HOME}}/bin/launcher.py run
stdout_logfile={{PRESTO_DATA_DIR}}/logs/presto.log
directory={{PRESTO_HOME}}/bin
redirect_stderr=true